# templates/nexus.service.epp
[Unit]
Description=Nexus Service
After=network.target

[Service]
Type=forking
LimitNOFILE=65536
ExecStart=<%= $nexus::install_path %>/nexus/bin/nexus start
ExecStop=<%= $nexus::install_path %>/nexus/bin/nexus stop
User=<%= $nexus::user %>
Restart=on-abort

[Install]
WantedBy=multi-user.target
